# Pypeline/supervisord.conf
; supervisord config for the gstreamer container.
;
; Five long-running processes — Xvfb, xcompmgr, overlay-server, Chromium,
; and the Python pipeline — are colocated in one container. Ordering is
; controlled by `priority` and per-program `startsecs`:
;
;   priority 1   Xvfb + xcompmgr   X11 virtual display + compositing manager.
;                                  Chromium's --enable-transparent-visuals
;                                  requires _NET_WM_CM_S0 from xcompmgr;
;                                  without it the body background renders
;                                  opaque even when CSS says transparent.
;   priority 2   overlay-server    Serves /overlays statically on :8080.
;   priority 3   chromium          Loads catch-counter.html in kiosk mode
;                                  against the Xvfb display.
;   priority 4   gstreamer         The Python pipeline. ximagesrc captures
;                                  whatever Chromium is rendering.
;
; All log streams are routed to supervisord's own stdout/stderr so that
; `docker logs gstreamer` shows everything in one place.
;
; The pipeline runs unconditionally; when browser_overlay.enabled=false in
; the YAML, the new branches are simply not instantiated. The browser stack
; (Xvfb/xcompmgr/Chromium/overlay-server) still runs — it's cheap when
; idle and avoids needing two image variants.

[supervisord]
nodaemon=true
user=root
logfile=/dev/null
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:xvfb]
command=Xvfb :99 -screen 0 1920x1080x24 -ac
autorestart=true
priority=1
startsecs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:overlay-server]
command=python3 /app/overlay_server.py --port 8080 --dir /overlays
autorestart=true
priority=2
startsecs=1
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:chromium]
; Wait for Xvfb's display socket before launching Chromium. Without the
; wrapper, Chromium races Xvfb and crashes "Can't open display" on first
; start; supervisord autorestarts, but that adds startup latency.
;
; The overlay HTML uses a solid #00ff00 background (chroma key green);
; BrowserOverlaySourceBranch's `alpha method=green` element keys that out,
; producing real per-pixel alpha for the compositor. No xcompmgr or
; transparent-visuals plumbing needed.
command=/bin/sh -c 'until xdpyinfo -display :99 >/dev/null 2>&1; do sleep 0.2; done; exec chromium --no-sandbox --disable-gpu --disable-dev-shm-usage --window-size=1920,1080 --window-position=0,0 --kiosk --app=http://localhost:8080/catch-counter.html'
environment=DISPLAY=":99"
autorestart=true
priority=3
startsecs=4
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:gstreamer]
command=python3 /app/main.py --config /app/configs/boat1.yaml
environment=DISPLAY=":99",PYTHONUNBUFFERED="1"
autorestart=true
priority=4
startsecs=6
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
